# docker/cron/crontab
* * * * * /home/app/src/bin/docma help > /home/app/logs/cron.log 2>&1
