# docker/cron/crontab
* * * * * /home/app/src/bin/docma build --source=/home/app/documents/ --destination=/home/app/results/ --logfile=/home/app/logs/docma.log > /home/app/logs/cron.log 2>&1
